# conf.d/novnc.conf
[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:%(ENV_VNC_PORT)s --listen 8080
autorestart=true